# scripts/systemd/creationDate-watcher.path.example
[Path]
Unit=creationDate-watcher.service
PathModified=<path to creationDate folder>

[Install]
WantedBy=default.target